# configs/profile.d/wsl.sh
#!/bin/sh

IS_WSL=`grep -i microsoft /proc/version`

export IS_WSL
